# herbstluftwm/services/battery/batteryCheck.service
[Unit]
Description=Battery health check

[Service]
ExecStart=%h/.config/herbstluftwm/scripts/battery/battery_health_check.sh
